# cron/crontab
MAILTO=""
49 * * * * ~/current/cron/wrapper.sh cronjob -E -s 'GMC Update' -r plu@pqpq.de -c ~/current/cron/update.pl
